# config/runit/netatalk/run
#!/bin/sh
sleep 10
exec netatalk -d -F /etc/afp.conf
