# config/runit/netatalk/run
#!/bin/sh
sleep 3
exec netatalk -d -F /etc/afp.conf
