# config/runit/netatalk/run
#!/bin/sh
sleep 5
exec netatalk -d -F /etc/afp.conf
