# config/runit/netatalk/run
#!/bin/sh
exec netatalk -d -F /etc/afp.conf
