# config/runit/avahi/run
#!/bin/sh
sleep 3
exec avahi-daemon --no-drop-root
