# config/runit/avahi/run
#!/bin/sh
exec avahi-daemon --no-drop-root
