# config/runit/avahi/run
#!/bin/sh
exec avahi-daemon --daemonize --no-drop-root
