# config/runit/dbus/run
#!/bin/sh
mkdir /var/run/dbus
exec dbus-daemon --system --nofork
